# clients/cpp-console/scripts/service/VoiceAssistant.timer
[Unit]
Description=Timer for the Voice Assistant

[Timer]
OnStartupSec=45

[Install]
WantedBy=default.target