# samples/clients/cpp-console/scripts/GGEC/HospitalityDemo.timer
[Unit]
Description=Timer for the hospitality demo

[Timer]
OnStartupSec=45

[Install]
WantedBy=default.target